# roles/sysstat/files/override.conf_sysstat-collect.timer
# Run system activity accounting tool every 2 minutes
[Timer]
OnCalendar=*:00/2

